# pkg-cache/pkg-cache.service
[Unit]
Description=Description of pkg-cache

[Service]
ExecStart=/usr/bin/pkg-cache.sh
Type=oneshot
User=root
