# aur-cache/aur-cache.service
[Unit]
Description=Description of aur-cache

[Service]
ExecStart=/usr/bin/aur-cache.sh
Type=oneshot
User=keishis
